# scripts/com.nachovuela.rastrillaje.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Agenda de NachoVuela para macOS (launchd).
  Corre el rastrillaje dos veces por día: 9:00 y 20:00.
  Instalación: ver scripts/agendar.md
  Acordate de reemplazar __RUTA__ por la ruta real del proyecto.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nachovuela.rastrillaje</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__RUTA__/scripts/run.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>__RUTA__/data/rastrillaje.log</string>
  <key>StandardErrorPath</key>
  <string>__RUTA__/data/rastrillaje.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
